# psd-frontend.service
[Unit]
Description=PSD Canvas Jaaz Frontend Service
After=network.target psd-backend.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cckz/psd-canvas-jaaz/react
Environment="NODE_ENV=production"
Environment="PORT=3004"
Environment="BACKEND_URL=http://127.0.0.1:57988"
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

